# archive/Linux/Processes/linux_lab_systemd_services/scripts/lab-demo.user.service
[Unit]
Description=Lab demo service for systemd practice
After=network.target

[Service]
Type=simple
ExecStart=REPLACE_WITH_LAB_ROOT/scripts/dummy_daemon.sh
Restart=no

[Install]
WantedBy=default.target
